# deploy/systemd/ai-tg-bot.service
[Unit]
Description=AI Telegram Bot
After=network.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/projects/ai-tg-bot
ExecStart=/usr/bin/node --env-file=/root/projects/ai-tg-bot/.env.bot src/bot/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-tg-bot
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
